# templates/default/logstash_agent.conf.erb
description "Logstash agent"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30
limit nofile 65550 65550

chdir <%= node['logstash']['server']['home'] %>

setuid <%= node['logstash']['user'] %>
<% unless node['logstash']['supervisor_gid'].to_s.empty? -%>
setgid <%= node['logstash']['supervisor_gid'] %>
<% end -%>

console log
script
  # Defaults
  PATH=/bin:/usr/bin
  LS_HOME=/var/lib/logstash
  LS_HEAP_SIZE="500m"
  LS_JAVA_OPTS="-Djava.io.tmpdir=${LS_HOME}"
  LS_LOG_FILE=/var/log/logstash/logstash.log
  LS_USE_GC_LOGGING=""
  LS_CONF_DIR=/etc/logstash/conf.d
  LS_OPEN_FILES=16384
  LS_NICE=19
  LS_OPTS=""

  # Override our defaults with user defaults:
  [ -f /etc/default/logstash_<%= @type %> ] && . /etc/default/logstash_<%= @type %>

  HOME="${HOME:-$LS_HOME}"
  JAVA_OPTS="${LS_JAVA_OPTS}"
  # Reset filehandle limit
  ulimit -n ${LS_OPEN_FILES}
  cd "${LS_HOME}"

  # Export variables
  export PATH HOME JAVA_OPTS LS_HEAP_SIZE LS_JAVA_OPTS LS_USE_GC_LOGGING
  test -n "${JAVACMD}" && export JAVACMD

  exec nice -n ${LS_NICE} ${HOME}/bin/logstash ${LS_OPTS}
end script

emits logstash-<%= @type %>-running
